# service/usr/lib/udev/rules.d/99-mm-rolling-linux-apps-port-types.rules
# do not edit this file, it will be overwritten on update
ACTION!="add|change|move|bind", GOTO="mm_rolling_linux_apps_port_types_end"
SUBSYSTEMS=="usb", ATTRS{bInterfaceNumber}=="?*", ENV{.MM_USBIFNUM}="$attr{bInterfaceNumber}", GOTO="mm_rw101_port_types"
GOTO="mm_rolling_linux_apps_port_types_end"

LABEL="mm_rw101_port_types"

ATTRS{idVendor}=="33f8", ATTRS{idProduct}=="0301",RUN+="/bin/sh -c '/sbin/modprobe usbserial && /sbin/modprobe option && echo 33f8 0301 ff > /sys/bus/usb-serial/drivers/option1/new_id'"

ATTRS{idVendor}=="33f8", ATTRS{idProduct}=="0302",RUN+="/bin/sh -c '/sbin/modprobe usbserial && /sbin/modprobe option && echo 33f8 0302 ff > /sys/bus/usb-serial/drivers/option1/new_id'"

ATTRS{idVendor}=="33f8", ATTRS{idProduct}=="01a8",RUN+="/bin/sh -c '/sbin/modprobe usbserial && /sbin/modprobe option && echo 33f8 01a8 ff > /sys/bus/usb-serial/drivers/option1/new_id'"

ATTRS{idVendor}=="33f8", ATTRS{idProduct}=="01a9",RUN+="/bin/sh -c '/sbin/modprobe usbserial && /sbin/modprobe option && echo 33f8 01a9 ff > /sys/bus/usb-serial/drivers/option1/new_id'"

ATTRS{idVendor}=="33f8", ATTRS{idProduct}=="0301", ENV{.MM_USBIFNUM}=="02", SUBSYSTEM=="tty", ENV{ID_MM_PORT_TYPE_AT_PRIMARY}="1"
ATTRS{idVendor}=="33f8", ATTRS{idProduct}=="0302", ENV{.MM_USBIFNUM}=="02", SUBSYSTEM=="tty", ENV{ID_MM_PORT_TYPE_AT_PRIMARY}="1"
ATTRS{idVendor}=="33f8", ATTRS{idProduct}=="0301", ENV{.MM_USBIFNUM}=="02", SUBSYSTEM=="tty", ENV{ID_MM_FIBOCOM_FASTBOOT}="1"
ATTRS{idVendor}=="33f8", ATTRS{idProduct}=="0301", ENV{.MM_USBIFNUM}=="02", SUBSYSTEM=="tty", ENV{ID_MM_PORT_TYPE_GPS}="0"
ATTRS{idVendor}=="33f8", ATTRS{idProduct}=="0302", ENV{.MM_USBIFNUM}=="02", SUBSYSTEM=="tty", ENV{ID_MM_FIBOCOM_FASTBOOT}="1"
ATTRS{idVendor}=="33f8", ATTRS{idProduct}=="01a8", ENV{.MM_USBIFNUM}=="03", SUBSYSTEM=="tty", ENV{ID_MM_FIBOCOM_FASTBOOT}="1"
ATTRS{idVendor}=="33f8", ATTRS{idProduct}=="01a9", ENV{.MM_USBIFNUM}=="03", SUBSYSTEM=="tty", ENV{ID_MM_FIBOCOM_FASTBOOT}="1"

LABEL="mm_rolling_linux_apps_port_types_end"

